# lsws-4.2.18/admin/misc/lshttpd.service.in
# should be added as /usr/lib/systemd/system/lshttpd.service
[Unit]
Description=The LiteSpeed HTTP Server
After=network.target remote-fs.target nss-lookup.target

[Service]
Type=forking
PIDFile=/tmp/lshttpd/lshttpd.pid
ExecStart=%LSWS_CTRL% start 
ExecReload=%LSWS_CTRL% restart 
ExecStop=%LSWS_CTRL% stop

KillMode=none
PrivateTmp=false

[Install]
WantedBy=multi-user.target
Alias=lsws
